# admin/configroot/etc/cron.d/tmbo
THEMAXX=/home/vagrant/sites/tmbo

# m     h       dom     mon     dow     cmd
# 4 0 * * * www-data $THEMAXX/offensive/zipOldFiles.pl 1 images/audio/ zips/audio

5 0 * * * www-data $THEMAXX/offensive/deleteOldFiles.pl 7 zips > /dev/null
2 0 * * * www-data $THEMAXX/offensive/zipYesterday.pl zips
3 1 * * * www-data $THEMAXX/offensive/deleteOldFiles.pl 2 quarantine > /dev/null
0 * * * * www-data /usr/bin/nice /usr/bin/php5 $THEMAXX/admin/commentIndexer.php